# dot_config/systemd/user/kbdlightd.service
[Unit]
Description=Auto keyboard backlight

[Service]
ExecStart=/usr/local/bin/kbdlightd
Restart=always

[Install]
WantedBy=default.target
